# tools/cc-channels/systemd/treadmill-channel@.service
[Unit]
Description=Treadmill orchestrator session (label %i)
After=default.target

[Service]
Type=simple
# Idempotent launcher: creates the tmux session if missing,
# then sends the Claude launch into it. Foreground (stays
# alive for systemd to supervise).
ExecStart=%h/treadmill/tools/cc-channels/systemd/treadmill-channel-launch %i
# Reap the label's whole session tree (tmux session + lock-holding
# launcher/claude) on every stop/restart AND after failed starts. The
# tree lives under the tmux server's cgroup, not this unit's, so
# KillMode can never reach it — this is the fix for the 2026-06-11
# orphan crashloop (NRestarts=4119; see tools/cc-channels AGENT.md).
ExecStopPost=%h/treadmill/tools/cc-channels/systemd/treadmill-channel-reap %i
# ADR-0091 §5 (task 92365367): an intentional stop delivers SIGTERM and
# the launcher exits 143 — without this, every scheduled pause lands the
# unit in 'failed' until reset-failed (papercut found 2026-06-12).
# OPERATOR STEP on rollout: the installed unit at
# ~/.config/systemd/user/treadmill-channel@.service is a plain copy —
# cp + `systemctl --user daemon-reload` to pick this up.
SuccessExitStatus=143
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
